# content/service/tachidesk/run
#!/bin/bash

# Dummy run for building image
if [ "${MODE}" = "build" ]; then
  sleep infinity
fi

# Run caddy
exec 2>&1
exec java -jar "/.tachiworkdir/tachidesk_latest.jar"
